# services/quiet-discord.service
[Unit]
Description=Quiet Discord listener (bridges DMs to Quiet)
After=network.target quiet-web.service
Requires=quiet-web.service

[Service]
Type=simple
WorkingDirectory=%h/quiet
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 %h/quiet/discord_listener.py --config %h/quiet/discord_config.json
Restart=always
RestartSec=10

# Logging — uses journal by default; override locally if you want file logs
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
